# etc/cron.d/agent
# Cadia agent - persistent job
* * *  * * import runjob --noduplicate --mail=rt-cadia@trellian.com --memory=16000 --tag=cadia:agent --warningtime=0 -- python3 /var/www/cadia.trillion.com/agent2.py --id=dns --scrapers=dns
# rdap now running on cadia01-san
#* * *  * * import runjob --noduplicate --mail=rt-cadia@trellian.com --memory=16000 --tag=cadia:agent --warningtime=0 -- python3 /var/www/cadia.trillion.com/agent.py --id=rdap-1 --scrapers=rdap --proxies=ipburger-1
#* * *  * * import runjob --noduplicate --mail=rt-cadia@trellian.com --memory=16000 --tag=cadia:agent --warningtime=0 -- python3 /var/www/cadia.trillion.com/agent.py --id=3 --scrapers=rdap,dns
#* * *  * * import runjob --noduplicate --mail=rt-cadia@trellian.com --memory=16000 --tag=cadia:agent --warningtime=0 -- python3 /var/www/cadia.trillion.com/agent.py --id=4 --scrapers=rdap,dns

# One agent dedicated for whois
# disabled, timing out, need to investigate
* * *  * * import runjob --noduplicate --mail=rt-cadia@trellian.com --memory=8000 --tag=cadia:agent --warningtime=0 -- python3 /var/www/cadia.trillion.com/agent2.py --id=whois --scrapers=whois --proxies=netnut
